# root/etc/services.d/Xorg/run
#!/bin/execlineb -P

# Redirect stderr to stdout.
fdmove -c 2 1

# Execute Xorg
/usr/bin/Xorg :1 -displayfd 3 -config /etc/X11/xrdp/xorg.conf -noreset -nolisten tcp
